# debian/openvswitch-brcompat.postinst
#!/bin/sh

set -e

# If openvswitch-switch is installed, and then later openvswitch-brcompat is
# installed, make sure that ovs-brcompatd starts.
if test X"$1" = Xconfigure && \
   test -x /etc/init.d/openvswitch-switch && \
   test -e /var/run/openvswitch/ovs-vswitchd.pid; then
    invoke-rc.d openvswitch-switch start || exit $?
fi

#DEBHELPER#

exit 0


