# rhel/usr_lib_systemd_system_openvswitch.service
[Unit]
Description=Open vSwitch
After=syslog.target network.target ovsdb-server.service
Requires=ovsdb-server.service

[Service]
Type=oneshot
ExecStart=/bin/true
ExecStop=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
