# xenserver/etc_init.d_openvswitch-xapi-update
#!/bin/bash
#
# openvswitch-xapi-update
#
# chkconfig: 2345 95 01
# description: Update Open vSwitch configuration from XAPI database at boot

# Copyright (C) 2009, 2010 Nicira, Inc.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at:
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.
### BEGIN INIT INFO
# Provides:      openvswitch-xapi-update
# Required-Start: $network $remote_fs
# Required-Stop: $network
# Default-Start:  3 5
# Default-Stop:
# Short-Description:    openvswitch-xapi-update
# Description:          reconfigures Open vSwitch based on XAPI configuration
### END INIT INFO

. /etc/init.d/functions

function do_host_call {
    xe host-call-plugin host-uuid="$INSTALLATION_UUID" plugin="openvswitch-cfg-update" fn="update" >/dev/null
}

function start {
    if [ ! -f /etc/xensource-inventory ]; then
        printf "openvswitch-xapi-update ERROR: XenSource inventory not present in /etc/xensource-inventory\n"
        exit 1
    fi

    if test -e /etc/xensource/network.conf; then
        NETWORK_MODE=$(cat /etc/xensource/network.conf)
    fi

    case ${NETWORK_MODE:=openvswitch} in
        vswitch|openvswitch)
            ;;
        bridge)
            exit 0
            ;;
        *)
            echo "Open vSwitch disabled (/etc/xensource/network.conf is invalid)" >&2
            exit 0
            ;;
    esac

    source /etc/xensource-inventory
    action "Updating configuration" do_host_call
}

case "$1" in
    start)
        start
        ;;
    stop)
        # Nothing to do here.
        ;;
    restart)
        start
        ;;
    help)
        printf "openvswitch-xapi-update [start|stop|restart]\n"
        ;;
    *)
        printf "Unknown command: $1\n"
        exit 1
        ;;
esac
